# deploy/systemd/rooney-dashboard.service
[Unit]
Description=Rooney Capital Trading Dashboard (Streamlit)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=/opt/pine/rooney-capital-v1/dashboard

# Environment configuration
Environment="PYTHONUNBUFFERED=1"

# Virtual environment - Run Streamlit
ExecStart=/opt/pine/rooney-capital-v1/venv/bin/streamlit run app.py --server.port 8501 --server.address 0.0.0.0 --server.headless true

# Restart configuration
Restart=on-failure
RestartSec=10
StartLimitInterval=350
StartLimitBurst=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rooney-dashboard

# Process limits
LimitNOFILE=4096

# Security hardening (optional)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
